# support/pressurecontrol.conf
start on (started mysql and net-device-up IFACE=eth0)
stop on runlevel [06]

expect fork
respawn

exec /usr/local/bin/pressurecontrol &

pre-stop exec /usr/bin/pkill -f pressurecontrol &
